# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php5-fpm]
command=/usr/sbin/php-fpm

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
